# contrib/docker.socket
[Socket]
SocketMode=0750
ListenStream=/var/run/docker.sock

[Install]
WantedBy=sockets.target
