# contrib/init/systemd/socket-activation/docker.socket
[Unit]
Description=Docker Socket for the API

[Socket]
ListenStream=/var/run/docker.sock

[Install]
WantedBy=sockets.target
