# puppet/modules/docker/templates/dockerd.conf
description     "Run dockerd"

stop on runlevel [!2345]
start on runlevel [3]

# if you want it to automatically restart if it crashes, leave the next line in
respawn

script
    test -f /etc/default/locale && . /etc/default/locale || true
    LANG=$LANG LC_ALL=$LANG /opt/go/bin/docker -d >> /var/log/dockerd 2>&1
end script
