# contrib/init/systemd/docker.service
[Unit]
Description=Docker Application Container Engine 
Documentation=http://docs.docker.io
Requires=network.target
After=multi-user.target

[Service]
Type=simple
ExecStartPre=/bin/mount --make-rprivate /
ExecStart=/usr/bin/docker -d

[Install]
WantedBy=multi-user.target
