# contrib/init/systemd/docker.service
[Unit]
Description=Docker Application Container Engine 
Documentation=http://docs.docker.io
After=network.target

[Service]
ExecStartPre=/bin/mount --make-rprivate /
ExecStart=/usr/bin/docker -d

[Install]
WantedBy=multi-user.target
