# packaging/ubuntu/docker.upstart
description     "Run docker"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]
respawn

script
    /usr/bin/docker -d
end script
